# setup/deb/toxphone/DEBIAN/prerm
#!/bin/bash

set +e

if [ ! -e /usr/lib/systemd/user/toxphone.service ]; then
    killall toxphone-watchdog 2> /dev/null
    sleep 1
    killall toxphone 2> /dev/null
fi

function execSuCommand()
{
    su -l $SUDO_USER -c "DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/$SUDO_UID/bus $1"
}

res=$(execSuCommand "systemctl --user list-unit-files | grep 'toxphone.service'")
if [ -n "$res" ]; then
    execSuCommand "systemctl --user stop    toxphone.service"
    execSuCommand "systemctl --user disable toxphone.service"
fi

res=$(systemctl list-unit-files | grep 'toxphone-cpufreq.service')
if [ -n "$res" ]; then
    systemctl stop    toxphone-cpufreq.service
    systemctl disable toxphone-cpufreq.service
fi

#rm -f /etc/xdg/autostart/toxphone-autorun.desktop

exit 0
